# etc/service/rocketchat/run
#!/bin/bash
# description "Rocketchat service manager"
sv start nginx
sv start rocketchat_mongo
chpst -u rocketchat:rocketchat exec /usr/local/bin/node /Rocket.Chat/bundle/main.js
