# etc/service/rocketchat/run
#!/bin/bash -e

exec chpst -u rocketchat:rocketchat /usr/local/bin/node /Rocket.Chat/bundle/main.js
